# mail-mta/dma/files/dma-flush.service
[Unit]
Description=Flush Dragonfly Mail Agent queue
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/dma -q
